# etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:activemq]
user=root
directory=/activemq
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=600
command=activemq start

